# daemon/rdp-display-manager.service
[Unit]
Description=RDP Display Switcher Daemon
Documentation=https://github.com/yourusername/rdp-display-switcher
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/share/rdp-display-switcher/rdp-display-manager.py
Restart=on-failure
RestartSec=5

# Import the full user session environment (needed for Wayland/KDE access)
ImportCredential=*
PassEnvironment=WAYLAND_DISPLAY XDG_RUNTIME_DIR DBUS_SESSION_BUS_ADDRESS

[Install]
WantedBy=graphical-session.target
